# misc/etc/rc.local
#!/bin/sh
#
# local startup daemons
#

echo -n 'starting local daemons:'

# Put local daemons here...

echo '.'
